# infra/vps/systemd/levelup-adspro-deploy.service
[Unit]
Description=LevelUP adsPRO Deploy Update
Wants=network-online.target
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/levelup-adspro
Environment=PROJECT_NAME=levelup-adspro
Environment=COMPOSE_FILE=docker-compose.vps.yml
ExecStart=/usr/bin/bash /opt/levelup-adspro/infra/vps/deploy.sh
