# infra/vps/systemd/levelup-adspro-deploy.timer
[Unit]
Description=Poll repo and deploy levelup-adspro changes

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
Unit=levelup-adspro-deploy.service
Persistent=true

[Install]
WantedBy=timers.target
